# cmd/deploy/cmd.conf
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 15 5
console log
script
	exec /bin/bash -lc "source /etc/cmd.env; cmd"
end script
