# packer/target-machine/etc/cron.hourly/refresh-authorized-keys
#!/bin/sh
/home/ubuntu/refresh-authorized-keys
