# packer/target-machine/etc/init/chain.conf
respawn
setuid ubuntu
setgid ubuntu
console log
start on runlevel [2345]
stop on runlevel [!2345]
